# files/ecr.ini
[supervisord]
nodaemon=true

[program:nginx]
command = nginx -g "daemon off;"
user = root
autostart = true

[program:crond]
command = crond -f
user = root
autostart = true
